# vlmcsd/files/vlmcsd.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/vlmcsd -i /etc/vlmcsd.ini -j /etc/vlmcsd.kmd -D
	procd_close_instance
}
